# air-quality-monitor.service
[Unit]
Description=Air Quality Monitor Service
After=multi-user.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi_air
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/pi/pi_air/air_quality_monitor.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/pi_air /var/log

[Install]
WantedBy=multi-user.target