# services/air-quality-monitor.service
[Unit]
Description=Air Quality Monitor Service
After=multi-user.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi_air
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/pi/pi_air/src/air_quality_monitor.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Reduced security hardening for compatibility
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target