# systemd/air-quality-monitor.service
[Unit]
Description=Air Quality Monitor for PMS7003 Sensor
After=network.target
Wants=network.target

[Service]
Type=simple
User={{USER}}
Group={{GROUP}}
WorkingDirectory={{PROJECT_DIR}}
Environment=PATH={{PROJECT_DIR}}/venv/bin
ExecStart={{PROJECT_DIR}}/venv/bin/python {{PROJECT_DIR}}/src/air_quality_monitor.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=air-quality-monitor

[Install]
WantedBy=multi-user.target